# tests/RedisConfigs/Docker.sentinel/supervisord.conf
[supervisord]
nodaemon=false

[program:redis-7010]
command=/usr/local/bin/redis-server /data/Sentinel/redis-7010.conf
directory=/data/Sentinel
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:redis-7011]
command=/usr/local/bin/redis-server /data/Sentinel/redis-7011.conf
directory=/data/Sentinel
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:sentinel-26379]
command=/usr/local/bin/redis-server /data/Sentinel/sentinel-26379.conf --sentinel
directory=/data/Sentinel
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:sentinel-26380]
command=/usr/local/bin/redis-server /data/Sentinel/sentinel-26380.conf --sentinel
directory=/data/Sentinel
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:sentinel-26381]
command=/usr/local/bin/redis-server /data/Sentinel/sentinel-26381.conf --sentinel
directory=/data/Sentinel
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true