# tests/RedisConfigs/Docker.cluster/supervisord.conf
[supervisord]
nodaemon=false

[program:cluster-7000]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7000.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:cluster-7001]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7001.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:cluster-7002]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7002.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:cluster-7003]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7003.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:cluster-7004]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7004.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:cluster-7005]
command=/usr/local/bin/redis-server /data/Cluster/cluster-7005.conf
directory=/data/Cluster
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
